# sh/profile.d/browser.sh
# Browser
if [ -n "$DISPLAY" ]
then
    BROWSER='firefox'
else
    BROWSER='lynx'
fi
export BROWSER

