# sh/profile.d/browser.sh
# Browser
if [ -n "$DISPLAY" ]; then
    BROWSER='google-chrome'
else
    BROWSER='lynx'
fi
export BROWSER

